# api/upes-api.service
[Unit]
Description=UPES-ECS Emergency PBX local status/control API (FastAPI/uvicorn on :8090)
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
User=root
# Run the app directly; it calls uvicorn.run() on 0.0.0.0:8090.
ExecStart=/usr/bin/python3 /opt/upes-ecs/api/upes_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
